# crontab
# Synchronize Google Drive
*/10 * * * * grive -p $HOME/GoogleDrive/ > /dev/null 2>&1
